# launchd/com.esg-marketplace.intake-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ESG Skills Marketplace — intake HTTP server LaunchAgent
  Keeps the intake_server.py running persistently.
  Cloudflare Tunnel exposes it at the URL configured in proxy.json.

  BEFORE LOADING:
    1. Generate a shared secret: python3 -c "import secrets; print(secrets.token_urlsafe(32))"
    2. Set that secret as MAC_MINI_INTAKE_SECRET below AND in your here.now account variable.
    3. Start a named Cloudflare Tunnel pointing to http://localhost:8765 (see docs/INTAKE-SETUP.md).
    4. Update site/.herenow/proxy.json upstream URL with your tunnel hostname.
    5. Copy this file to ~/Library/LaunchAgents/ and load it:
         cp launchd/com.esg-marketplace.intake-server.plist ~/Library/LaunchAgents/
         launchctl load ~/Library/LaunchAgents/com.esg-marketplace.intake-server.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.esg-marketplace.intake-server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/ricardofilho/Documents/Projects/active/esg-skills-marketplace/scripts/intake_server.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ricardofilho/Documents/Projects/active/esg-skills-marketplace</string>

  <!-- Replace placeholder values before loading -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>MAC_MINI_INTAKE_SECRET</key>
    <string>REPLACE_WITH_SHARED_SECRET</string>
    <key>INTAKE_PORT</key>
    <string>8765</string>
    <key>ESG_INTAKE_ROOT</key>
    <string>/Users/ricardofilho/Documents/Projects/active/esg-skills-marketplace</string>
  </dict>

  <!-- Restart automatically if the server crashes -->
  <key>KeepAlive</key>
  <true/>

  <!-- Bounded log files -->
  <key>StandardOutPath</key>
  <string>/Users/ricardofilho/Library/Logs/esg-marketplace/intake-server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ricardofilho/Library/Logs/esg-marketplace/intake-server-error.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
